# bin/news-feed/templates/startnewsfeed.service
[Unit]
Description=Newsfeed Service
Documentation=
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/newsfeed
Environment="APP_PORT=8082"
ExecStart=/usr/bin/bash /opt/newsfeed/startnewsfeed.sh >/dev/null 2>&1
Restart=on-failure

[Install]
WantedBy=multi-user.target
